# birdnetpi/birdnet_audio_bridge.service
[Unit]
Description=ESP32 UDP Audio to ALSA Loopback
After=network.target sound.target

[Service]
ExecStart=/opt/birdnet-audio-bridge/birdnet_audio_bridge.sh
Restart=always
RestartSec=5
User=jan

[Install]
WantedBy=multi-user.target